# roles/asreview/files/asreview.service
[Unit]
Description=Starting ASReview Lab when a user logs in.

[Service]
ExecStart=asreview lab

[Install]
WantedBy=multi-user.target